# deadline-collector.service
# Systemd service file for Deadline Collector
# Installation instructions:
#
# 1. Edit this file and update paths:
#    - WorkingDirectory: /path/to/deadline-collector
#    - ExecStart: /path/to/deadline-collector/venv/bin/python
#    - User: your-username
#
# 2. Copy to systemd directory:
#    sudo cp deadline-collector.service /etc/systemd/system/
#
# 3. Reload systemd:
#    sudo systemctl daemon-reload
#
# 4. Enable service to start on boot:
#    sudo systemctl enable deadline-collector
#
# 5. Start the service:
#    sudo systemctl start deadline-collector
#
# 6. Check status:
#    sudo systemctl status deadline-collector
#
# 7. View logs:
#    sudo journalctl -u deadline-collector -f

[Unit]
Description=Deadline Article Collector Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/path/to/deadline-collector
ExecStart=/path/to/deadline-collector/venv/bin/python collector.py --schedule

# Restart policy
Restart=always
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=deadline-collector

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
